# systemd/sailly-backend-api.service.template
[Unit]
Description=Sailly Backend API
After=network.target

[Service]
Type=simple
User=charles2
WorkingDirectory=/opt/sailly/backend-api
Environment=NODE_ENV=production
EnvironmentFile=/opt/sailly/backend-api.env
ExecStart=/usr/bin/node --require /opt/sailly/backend-api/node_modules/tsx/cjs.cjs /opt/sailly/backend-api/index-minimal.ts
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
