# voice-agent/sound-validation-watchdog.service
[Unit]
Description=Sailly Sound Validation — Autonomous STS Caller Test Runner
Documentation=https://github.com/sailly/sailly-browser-demo
After=network.target sailly-browser-demo.service
Wants=sailly-browser-demo.service

[Service]
Type=simple
User=charles2_hotmail_de
Group=charles2_hotmail_de
WorkingDirectory=/home/charles2/sailly-browser-demo

# Load all API keys from the .env file
EnvironmentFile=/home/charles2/sailly-browser-demo/.env

ExecStart=/home/charles2/sailly-browser-demo/run_validation.sh \
    --daemon \
    --phase all \
    --resume

# Restart on failure — this is the "autonomous" part.
# If the process crashes, systemd will restart it with --resume so it picks up
# from the last completed checkpoint instead of re-running everything.
Restart=on-failure
RestartSec=30s
StartLimitIntervalSec=300
StartLimitBurst=5

# Write stdout/stderr to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sound-validation

# Send SIGTERM first, wait 30s, then SIGKILL
TimeoutStopSec=30s

# Resource limits — avoid overwhelming the VM
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
