# voice-agent/deploy/builder_lab/sailly-builder-lab-dashboard.service.template
[Unit]
Description=Sailly Builder Lab Dashboard (Next.js, isolated)
After=network.target sailly-builder-lab.service
Wants=sailly-builder-lab.service

[Service]
Type=simple
User=charles2
WorkingDirectory=/home/charles2/builder-lab/sailly/apps/dashboard
Environment=PORT=13001
Environment=NODE_ENV=production
Environment=VOICE_AGENT_ORIGIN=http://127.0.0.1:18080
Environment=VOICE_AGENT_URL=http://127.0.0.1:18080
Environment=NEXT_PUBLIC_VOICE_AGENT_WSS=ws://127.0.0.1:18080
Environment=PATH=/home/charles2/.nvm/versions/node/v20.20.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStartPre=/bin/bash -c "cp -r /home/charles2/builder-lab/sailly/apps/dashboard/.next/static /home/charles2/builder-lab/sailly/apps/dashboard/.next/standalone/.next/static 2>/dev/null || true"
ExecStart=/home/charles2/.nvm/versions/node/v20.20.1/bin/node .next/standalone/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
